# confirm-alive.crontab
0 */12 * * * root cd /app && cargo run --bin confirm-alive >> /var/log/cron.log 2>&1
